# firmware/pi_robot/pi-robot-rtc.service
[Unit]
Description=Better Robotics — Pi WebRTC peer (aiortc)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Runs as the robot user — the PTY spawned for `shell` channels lands at
# this UID, matching `ssh robot@<host>.local`. Independent of pi-robot
# (which still runs as root for BLE), so a crash here doesn't touch
# motors / LEDs.
User=robot
WorkingDirectory=__HOME__/better-robotics/firmware/pi_robot
ExecStart=__HOME__/better-robotics/firmware/pi_robot/.venv/bin/python pi_robot_rtc.py
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
