# public/firmware/pi_robot/pi-robot.service
[Unit]
Description=Better Robotics — Pi robot firmware
After=bluetooth.service NetworkManager.service
Requires=bluetooth.service

[Service]
Type=simple
# Root is required to register BLE advertisements via BlueZ on Pi OS.
User=root
WorkingDirectory=__HOME__/better-robotics/firmware/pi_robot
ExecStart=__HOME__/better-robotics/firmware/pi_robot/.venv/bin/python pi_robot.py
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
