# public/firmware/pi_robot/usb-gadget.service
[Unit]
Description=USB composite gadget (ECM + ACM) for recovery access
# Must run before NetworkManager so usb0 is present for NM to configure.
# Must NOT depend on pi-robot — the whole point is being independent of it.
After=systemd-modules-load.service
Before=NetworkManager.service
ConditionPathExists=/sys/class/udc

[Service]
Type=oneshot
ExecStart=/usr/local/bin/usb-gadget-setup.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
